# etc/services.d/stunnel/run
#!/usr/bin/execlineb -P

/usr/bin/stunnel /etc/stunnel/stunnel.conf
